# systemd/claude-performance-monitor.timer
[Unit]
Description=Weekly Claude Code hooks performance report generation

[Timer]
# Run every Sunday at 20:00
OnCalendar=Sun *-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
