# systemd/claude-extended-mcp.service
[Unit]
Description=Claude Extended MCP Server Launcher
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/claude-extended-launcher.py daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin
Environment=HOME=%h

# Security (optional hardening)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target