# systemd/claude-state-validator.service
[Unit]
Description=Validate and backup Claude Code state files
Documentation=man:claude-state-validator(1)
After=network.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/claude-state-validator.sh
MemoryMax=50M
CPUQuota=10%%
StandardOutput=journal
StandardError=journal
TimeoutSec=60
Restart=no

[Install]
WantedBy=default.target
